# unix/os/linux/arch/systemd/system/getty@tty1.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- jon' --noclear --skip-login - $TERM
