# unix/os/linux/arch/getty@tty1.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- joni' --noclear --skip-login - $TERM
