# unix/os/linux/etc/systemd/system/backup-backblaze-b2.timer
[Unit]
Description=back up of home directory to Backblaze B2 daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

